# systemd/user/llama-server.service
[Unit]
Description=Llama.cpp AI Model Server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=HSA_OVERRIDE_GFX_VERSION=11.5.1
Environment=ROCBLAS_USE_HIPBLASLT=1
EnvironmentFile=%h/.config/llama-server/current.env
ExecStart=/usr/bin/bash -lc 'exec /usr/bin/llama-server ${LLAMA_MODEL_ARGS} ${LLAMA_SERVER_ARGS}'
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=default.target
